# docs/goauthing.timer
[Unit]

[Timer]
OnBootSec=5
OnUnitActiveSec=1min
#OnCalendar=hourly

[Install]
WantedBy=timers.target

